# root/etc/hotplug.d/iface/099-ustb
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

/usr/share/USTB/ipv6.sh
