# root/etc/hotplug.d/iface/099-ustb
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

nat=$(uci get ustb.@advance[0].nat)
if [ "$nat" != "1" ]; then
	return
fi

sleep 5

/usr/sare/ustb/ipv6.sh
