# root/etc/init.d/ustb
#!/bin/sh /etc/rc.common
START=50

start() {
	echo -n "" >/usr/share/USTB/link.log
	rm -rf /tmp/ustb.log
	enable=$(uci get ustb.@login[0].enable)
	if [ "$enable" == 1 ]; then
		echo "[$(date +'%Y-%m-%d %H:%M:%S')] Service start" >>/usr/share/USTB/link.log
		/usr/share/USTB/ipv6.sh
		/usr/share/USTB/link.sh &
	fi
}

stop() {
	pid=$(ps | grep link.sh | grep -v grep | awk '{print $1}')
	if [ "$pid" != "" ]; then
		kill -9 $pid
	fi
	echo -n "" >/usr/share/USTB/link.log
	rm -rf /tmp/ustb.log
}
